# cron/crontabs/root
0 0 * * * echo 'HELLO' 
0 0 * * * curl -X POST 'http://vpn_bot:6060/api' 