# cron/crontabs/root
* * * * * echo 'HELLO' 
* * * * * curl -X GET 'http://web_server:6000/execute-script'