# cron/crontabs/root
* * * * * echo 'HELLO' 
* * * * * curl -X POST 'http://vpn_bot:6060/api' 